# chef_repo/cookbooks/my_lamp/files/default/httpd.conf
start on vagrant-mounted
stop on stopping network-services

exec /sbin/service httpd start